# supervisord.conf
[supervisord]
nodaemon=true

[program:query-pdb-server]
command=/query-pdb/build/server/query_pdb_server --log=1
autostart=true
autorestart=true
